# install/systemd/user/arcade-launcher.service
[Unit]
Description=c-base Arcade Launcher (session user service)
# Starts with the graphical session, inherits its display environment,
# stops when the session ends.
After=graphical-session.target
PartOf=graphical-session.target

[Service]
# CWD is the data root: the launcher resolves games/, scores/, logs/ under it.
# The binary sits in /arcade itself because exported Godot chdirs to the exe
# dir at startup — exe dir and WorkingDirectory must both be the data root.
WorkingDirectory=/arcade
ExecStart=/arcade/launcher.x86_64
Restart=always
RestartSec=2

[Install]
WantedBy=graphical-session.target
